# 80-realsense.rules
SUBSYSTEM=="video4linux", ATTRS{idVendor}=="8086", ATTRS{idProduct}=="0a66",RUN+="/usr/local/bin/realsense-map-controls '%E{DEVNAME}'"

